# petclinic.service
[Unit]
Description=petclinic Service

[Service]
ExecStart=/opt/jars/spring-petclinic-2.7.0-SNAPSHOT.jar 

[Install]
WantedBy=multi-user.target